# systemd/kde-lockscreen-refill.timer
[Unit]
Description=Periodic Picsum refill for the KDE lockscreen cache

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=false

[Install]
WantedBy=timers.target
