# systemd/kde-lockscreen-refill.service
[Unit]
Description=KDE lockscreen Picsum refill (on-demand)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/lib/kde-lockscreen/venv/bin/python -m kde_lockscreen.fetcher --mode refill
Nice=15
